# deploy/xlstatus-agent.service
[Unit]
Description=XLStatus Agent
After=network.target

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/opt/xlstatus-agent
ExecStart=/usr/local/bin/xlstatus-agent run --config /etc/xlstatus-agent/agent.json
Restart=always
RestartSec=5s

# Environment
Environment="RUST_LOG=info"

# Security
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/xlstatus-agent /etc/xlstatus-agent

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=xlstatus-agent

[Install]
WantedBy=multi-user.target
