# deploy/xlstatus-server.service
[Unit]
Description=XLStatus Server
After=network.target

[Service]
Type=simple
User=xlstatus
Group=xlstatus
WorkingDirectory=/opt/xlstatus
ExecStart=/usr/local/bin/xlstatus-server
Restart=on-failure
RestartSec=5s

# Environment
Environment="CONFIG_FILE=/etc/xlstatus/server.toml"
Environment="RUST_LOG=info"
Environment="HTTP_BIND=0.0.0.0:8080"
Environment="GRPC_BIND=0.0.0.0:50051"

# Security
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/xlstatus

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=xlstatus-server

[Install]
WantedBy=multi-user.target
